# systemd/sanlight-meshd-generic.service.example
[Unit]
Description=SANlight BlueZ Mesh Daemon using generic HCI
Documentation=https://github.com/Nibbels/sanlight-mesh-bluez-poc
After=network-online.target
Wants=network-online.target
Conflicts=bluetooth.service bluetooth-mesh.service

[Service]
Type=simple
# Keep the unit itself minimal. Controller cleanup is done by scripts/install-service.sh
# before starting the service. Some Bluetooth helper commands can block when used
# as ExecStartPre under systemd.
ExecStart=/usr/libexec/bluetooth/bluetooth-meshd --io generic:hci0 --nodetach
Restart=on-failure
RestartSec=3
TimeoutStartSec=20

[Install]
WantedBy=multi-user.target
