# systemd/sanlight-meshd-generic.service.example
[Unit]
Description=SANlight BlueZ Mesh daemon using generic HCI
Documentation=https://github.com/Nibbels/sanlight-mesh-bluez-poc
After=dbus.service sys-subsystem-bluetooth-devices-hci0.device
Wants=dbus.service sys-subsystem-bluetooth-devices-hci0.device
Conflicts=bluetooth.service bluetooth-mesh.service
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=dbus
BusName=org.bluez.mesh
ExecStart=/usr/local/libexec/sanlight-meshd-generic-start
Restart=on-failure
RestartSec=3
TimeoutStartSec=25
TimeoutStopSec=10
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/bluetooth/mesh

[Install]
WantedBy=multi-user.target
