# systemd/sanlight-meshd-generic.service.example
[Unit]
Description=SANlight BlueZ Mesh Daemon using generic HCI
Documentation=https://github.com/Nibbels/sanlight-mesh-bluez-poc
After=network-online.target
Wants=network-online.target
Conflicts=bluetooth.service bluetooth-mesh.service

[Service]
Type=simple
ExecStartPre=/usr/bin/rfkill unblock bluetooth
ExecStartPre=-/usr/sbin/hciconfig hci0 down
ExecStartPre=-/usr/bin/btmgmt --index 0 power off
ExecStart=/usr/libexec/bluetooth/bluetooth-meshd --io generic:hci0 --nodetach
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
